# supervisor.ini
[program:VisitArticle]
directory=/src/yuehu/yuehu
command=celery -A celery_yuehu worker --loglevel=error --concurrency=1 -n queue-article-visit -Q yuehu-article-visit
user=root
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/src/yuehu/logs/worker-error-visit-article.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB

[program:LikeArticle]
directory=/src/yuehu/yuehu
command=celery -A celery_yuehu worker --loglevel=error --concurrency=1 -n queue-article-like -Q yuehu-article-like
user=root
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/src/yuehu/logs/worker-error-like-article.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB

[program:CommentArticle]
directory=/src/yuehu/yuehu
command=celery -A celery_yuehu worker --loglevel=error --concurrency=1 -n queue-article-comment -Q yuehu-article-comment
user=root
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/src/yuehu/logs/worker-error-comment-article.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
